# backend/crontab
0 0 * * * /app/daily.sh >> /app/src/workspace/cron.log 2>&1

# need EOF